# crontab
0 1 * * * bash /opt/sbr/dl/auto-download.sh
